# systemd-boot/.config/systemd/user/firefox.service
[Unit]
Description=Firefox
Requires=network-online.target

[Service]
ExecStart=firefox
