# systemd-boot/.config/systemd/user/vesktop.service
[Unit]
Description=Vesktop
Requires=network-online.target

[Service]
ExecStart=vesktop
