# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter

[Service]
ExecStart=udiskie

[Install]
WantedBy=default.target
